# scripts/forge.service
[Unit]
Description=Forge Backend Platform
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/home/ubuntu/forge serve --config forge.yaml
Restart=always
RestartSec=5
LimitNOFILE=65536
Environment="FORGE_ENV=production"

[Install]
WantedBy=multi-user.target
